# linux/touchpad-dwt/touchpad-dwt.service
[Unit]
Description=Disable touchpad while typing (MacBook palm-rejection workaround)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/touchpad-dwt
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
